# deployment/ciris-manager.service
[Unit]
Description=CIRIS Manager - Agent Lifecycle Management Service
After=docker.service
Requires=docker.service

[Service]
Type=exec
User=ciris-manager
Group=ciris
WorkingDirectory=/opt/ciris-manager
EnvironmentFile=/etc/ciris-manager/environment
ExecStart=/opt/ciris-manager/venv/bin/ciris-manager --config /etc/ciris-manager/config.yml
Restart=always
RestartSec=10

# Docker socket access
SupplementaryGroups=docker

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ciris-manager

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/ciris/agents /var/log/ciris-manager /etc/ciris-manager /home/ciris/nginx
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
LockPersonality=true

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target