# deployment/ciris-backup.timer
[Unit]
Description=Daily backup of CIRIS Manager

[Timer]
# Run daily at 2 AM
OnCalendar=*-*-* 02:00:00
# Randomize by up to 30 minutes to avoid load spikes
RandomizedDelaySec=1800
# Run immediately if system was off during scheduled time
Persistent=true

[Install]
WantedBy=timers.target